# contrib/systemd/open-sesame-wm.service
[Unit]
Description=Open Sesame window manager daemon
Documentation=https://github.com/scopecreep-zip/open-sesame
Requires=open-sesame-profile.service
After=open-sesame-profile.service
PartOf=open-sesame.target

[Service]
Type=notify
ExecStart=%h/.cargo/bin/daemon-wm
Restart=on-failure
RestartSec=5
WatchdogSec=30
NoNewPrivileges=yes

# Security hardening
PrivateNetwork=no
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=%t/open-sesame
ReadWritePaths=%h/.cache/open-sesame

# Resource limits
LimitNOFILE=4096
MemoryMax=128M

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=open-sesame.target
